# templates/memory-wiki-audit.service.template
[Unit]
Description=Cairnkeep memory + wiki invalidation audit

[Service]
Type=oneshot
ExecStart=@@AUDIT_SCRIPT@@ --para-root @@PARA_ROOT@@ --report @@REPORT_PATH@@
